# root/default_skel/etc/inittab
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys

#::sysinit:/bin/mount -t tmpfs -o size=64k,mode=0755 tmpfs /dev
#::sysinit:/bin/mkdir /dev/pts
::sysinit:/bin/mount -t devpts devpts /dev/pts
::sysinit:/bin/echo /sbin/mdev > /proc/sys/kernel/hotplug
#::sysinit:/sbin/sysctl -w kernel.hotplug=/sbin/mdev
::sysinit:/sbin/mdev -s

::sysinit:/sbin/hwclock -s
::sysinit:/bin/hostname -F /etc/hostname

::sysinit:/sbin/syslogd
::sysinit:/sbin/klogd

::sysinit:/etc/rc

::askfirst:-/bin/sh
#::respawn:/sbin/getty 115200 ttyS0

::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a
