# root/default_skel/etc/inittab
::sysinit:/etc/rc

# Start an "askfirst" shell on the console (whatever that may be)
::askfirst:-/bin/sh

#::respawn:/sbin/getty 115200 ttyS0

# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a

