# root/default_skel/etc/inittab
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/mount -t devpts devpts /dev/pts

::sysinit:/sbin/hwclock -f /dev/rtc -s
::sysinit:/bin/hostname -F /etc/hostname

::sysinit:/sbin/syslogd
::sysinit:/sbin/klogd

::sysinit:/etc/rc

::askfirst:-/bin/sh
#::respawn:/sbin/getty 115200 ttyS0

::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a
